# deploy/townlet-demo.service
[Unit]
Description=Townlet Multi-Day Demo Training
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKDIR%
ExecStart=%VENV_PYTHON% -m townlet.demo.runner %CONFIG% %DB_PATH% %CHECKPOINT_DIR%
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (will be substituted by install script)
Environment="TOWNLET_CONFIG=%CONFIG%"
Environment="TOWNLET_DB_PATH=%DB_PATH%"
Environment="TOWNLET_CHECKPOINT_DIR=%CHECKPOINT_DIR%"

[Install]
WantedBy=multi-user.target
